# manager/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

# ============================================================
# Infrastructure Layer (start first)
# ============================================================

[program:minio]
command=/opt/hiclaw/scripts/init/start-minio.sh
priority=50
autorestart=true
stdout_logfile=/var/log/hiclaw/minio.log
stderr_logfile=/var/log/hiclaw/minio-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

[program:tuwunel]
command=/opt/hiclaw/scripts/init/start-tuwunel.sh
priority=100
autorestart=true
stdout_logfile=/var/log/hiclaw/tuwunel.log
stderr_logfile=/var/log/hiclaw/tuwunel-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

# ============================================================
# Higress Gateway Layer
# Uses original Higress start scripts at /usr/local/bin/ with
# environment variables from the official Higress supervisord.conf
# ============================================================

[program:higress-apiserver]
directory=/
command=bash /usr/local/bin/start-apiserver.sh
priority=200
autorestart=true
startsecs=1
stdout_logfile=/var/log/hiclaw/higress-apiserver.log
redirect_stderr=true
stdout_logfile_maxbytes=10MB
environment=ENABLE_CONSOLE_ROUTE="1"

[program:higress-controller]
directory=/
command=bash /usr/local/bin/start-controller.sh
priority=300
autorestart=true
startsecs=1
stdout_logfile=/var/log/hiclaw/higress-controller.log
redirect_stderr=true
stdout_logfile_maxbytes=10MB
environment=POD_NAME="higress-controller",POD_NAMESPACE="higress-system",CONTROLLER_KEEP_XDS_CONFIG_LABELS="false",CONTROLLER_KEEP_XDS_CONFIG_ANNOTATIONS="false",PILOT_ENABLE_GATEWAY_API="false",PILOT_ENABLE_ALPHA_GATEWAY_API="false",ENABLE_LEADER_ELECTION="false"

[program:higress-pilot]
directory=/
command=bash /usr/local/bin/start-pilot.sh
priority=400
autorestart=true
startsecs=1
stdout_logfile=/var/log/hiclaw/higress-pilot.log
redirect_stderr=true
stdout_logfile_maxbytes=10MB
environment=REVISION="default",JWT_POLICY="none",POD_NAME="higress-pilot",POD_NAMESPACE="higress-system",ROOT_CA_DIR="/etc/certs",PILOT_CERT_PROVIDER="istiod",KUBECONFIG="/app/kubeconfig",PILOT_ENABLE_LDS_CACHE="false",PILOT_ENABLE_CROSS_CLUSTER_WORKLOAD_ENTRY="false",PILOT_ENABLE_METADATA_EXCHANGE="false",PILOT_SCOPE_GATEWAY_TO_NAMESPACE="true",VALIDATION_ENABLED="false",PILOT_TRACE_SAMPLING="1",PILOT_ENABLE_PROTOCOL_SNIFFING_FOR_OUTBOUND="true",PILOT_ENABLE_PROTOCOL_SNIFFING_FOR_INBOUND="true",PILOT_ENABLE_ANALYSIS="false",PILOT_ENABLE_GATEWAY_API="false",PILOT_ENABLE_ALPHA_GATEWAY_API="false",PILOT_ENABLE_GATEWAY_API_STATUS="false",PILOT_ENABLE_GATEWAY_API_DEPLOYMENT_CONTROLLER="false",PILOT_ENABLE_ALPN_FILTER="true",VALIDATION_WEBHOOK_CONFIG_NAME="",ISTIO_DUAL_STACK="false",ENABLE_OPTIMIZED_CONFIG_REBUILD="false",PILOT_ENABLE_K8S_SELECT_WORKLOAD_ENTRIES="false",DEFAULT_UPSTREAM_CONCURRENCY_THRESHOLD="10000",ISTIO_GPRC_MAXRECVMSGSIZE="104857600",ENBALE_SCOPED_RDS="true",ON_DEMAND_RDS="false",HOST_RDS_MERGE_SUBSET="false",ENABLE_LEADER_ELECTION="false",PRIORITIZED_LEADER_ELECTION="false",INJECT_ENABLED="false",CLUSTER_ID="Kubernetes",CUSTOM_CA_CERT_NAME="higress-ca-root-cert",DEBUG_AUTH="false"

[program:higress-gateway]
directory=/
command=bash /usr/local/bin/start-gateway.sh
priority=500
autorestart=true
startsecs=1
stdout_logfile=/var/log/hiclaw/higress-gateway.log
redirect_stderr=true
stdout_logfile_maxbytes=10MB
environment=JWT_POLICY="none",CA_ROOT_CA="/etc/certs/root-cert.pem",XDS_ROOT_CA="/etc/certs/root-cert.pem",PROV_CERT="/etc/certs",PILOT_XDS_SEND_TIMEOUT="60s",PROXY_XDS_VIA_AGENT="true",ENABLE_INGRESS_GATEWAY_SDS="false",POD_NAME="higress-gateway",POD_NAMESPACE="higress-system",ISTIO_META_CLUSTER_ID="Kubernetes",ISTIO_META_HTTP10="1",INSTANCE_NAME="higress-gateway"

[program:higress-console]
directory=/
command=bash /usr/local/bin/start-console.sh
priority=600
autorestart=true
startsecs=1
stdout_logfile=/var/log/hiclaw/higress-console.log
redirect_stderr=true
stdout_logfile_maxbytes=10MB

# ============================================================
# Element Web (Nginx static file server)
# ============================================================

[program:element-web]
command=/opt/hiclaw/scripts/init/start-element-web.sh
priority=650
autorestart=true
stdout_logfile=/var/log/hiclaw/element-web.log
stderr_logfile=/var/log/hiclaw/element-web-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

# ============================================================
# Application Layer
# ============================================================

[program:setup-host-symlinks]
command=/opt/hiclaw/scripts/setup-host-symlinks.sh
priority=750
autostart=true
autorestart=false
startretries=0
startsecs=0
stdout_logfile=/var/log/hiclaw/setup-host-symlinks.log
stderr_logfile=/var/log/hiclaw/setup-host-symlinks-error.log
stdout_logfile_maxbytes=1MB

[program:mc-mirror]
command=/opt/hiclaw/scripts/init/start-mc-mirror.sh
priority=700
autorestart=true
stdout_logfile=/var/log/hiclaw/mc-mirror.log
stderr_logfile=/var/log/hiclaw/mc-mirror-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

[program:manager-agent]
command=/opt/hiclaw/scripts/init/start-manager-agent.sh
priority=800
autorestart=true
stdout_logfile=/var/log/hiclaw/manager-agent.log
stderr_logfile=/var/log/hiclaw/manager-agent-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

[program:manager-watch]
command=/opt/hiclaw/scripts/init/start-manager-watch.sh
priority=900
autorestart=true
stdout_logfile=/var/log/hiclaw/manager-watch.log
stderr_logfile=/var/log/hiclaw/manager-watch-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
